# scripts/maintenance/com.csp.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Claude Science Proxy (CSP) daily maintenance launchd agent.
  Fires at 09:00 and 21:00 local time (Asia/Shanghai in daily-maintenance.sh).
  __REPO__ is substituted by install-maintenance.sh on install.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.csp.maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/scripts/maintenance/daily-maintenance.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__REPO__/findings/auto-maint/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/findings/auto-maint/logs/launchd.err.log</string>
</dict>
</plist>
