# scripts/com.csswitch.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  CSSwitch 每日维护巡检 launchd agent。
  每天本地时间 09:00 与 21:00（Asia/Shanghai）触发 daily-maintenance.sh。
  绝对路径写死本机；仓库里这份是权威副本，install-maintenance.sh 会拷到 ~/Library/LaunchAgents/。
  launchd 若在触发时机器休眠，会在唤醒后补跑一次（合并错过的日历事件）。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.csswitch.maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/superjj/ccproj/CSswitch/scripts/daily-maintenance.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/superjj/ccproj/CSswitch/findings/auto-maint/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/superjj/ccproj/CSswitch/findings/auto-maint/logs/launchd.err.log</string>
</dict>
</plist>
